# supervisord.conf
[inet_http_server]
port=127.0.0.1:9001
username=Bill
password= TestingConfig


[supervisorctl]
serverurl=http://127.0.0.1:8000
username= flaskgram


[program:foo]
command = /bin/cat