# services.d/isso.ini
[program:isso]
command = %(ENV_HOME)s/.local/bin/isso -c %(ENV_HOME)s/user.cfg run
autostart = true
autorestart = true